# voice/stt_pipeline/claudette-stt.service
[Unit]
Description=Claudette Home — STT Transcription API (Whisper)
After=network.target
Documentation=https://github.com/mattcarp/mc-home

[Service]
Type=simple
User=sysop
WorkingDirectory=/home/sysop/projects/mc-home/voice/stt_pipeline
ExecStart=/usr/bin/python3 -m uvicorn transcribe_api:app --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=5
EnvironmentFile=/etc/environment
Environment=WHISPER_MODEL=base.en
Environment=WHISPER_LANGUAGE=en
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp
PrivateTmp=true

[Install]
WantedBy=multi-user.target
